# .config/sv/xkeysnail/run
#!/usr/bin/env bash

nq bash -c "sleep 3; ~/.scripts/xsettings.bash"
exec xkeysnail
