# .config/sv/pipewire/run
#!/usr/bin/env bash

pipewire & 
sleep 1
wireplumber & 
sleep 1
nq bash -c "sleep 3; ~/.scripts/xsettings.bash"
exec pipewire-pulse
