# .config/sv/pipewire/run
#!/usr/bin/env bash

nq bash -c "sleep 1; sv restart ../wireplumber"
exec nice -n -15 pipewire
