# .config/sv/keynav/run
#!/usr/bin/env bash

# export DISPLAY=:0
exec keynav
