# .config/sv/xkeysnail-fix-after/run
#!/usr/bin/env bash

sleep 5
exec entr -n ~/.scripts/xsettings.bash < <(echo "${HOME}/.local/share/service-graphic/xkeysnail/supervise/status")
