# .config/sv/sxhkd/run
#!/usr/bin/env bash

cd
SHELL=sh exec sxhkd
