# .config/sv/detox/run
#!/usr/bin/env bash

while true; do
	detox ~/Pictures
	sleep 2h
done
