# .config/sv/polybar/run
#!/usr/bin/env bash

while [ ! -r /run/user/1000/pulse/native ]; do
	sleep 0.5
done
exec polybar
