# .config/sv/polybar/run
#!/usr/bin/env bash

[ -r /run/user/1000/pulse/native ] || exit 1
cd ~/.config/polybar/
m4 config.m4 > config.ini
cd
exec polybar
