# .config/sv/wallpapers/run
#!/usr/bin/env bash

exec ~/.scripts/Wallpapers.sh
