# .config/sv/wallpapers/run
#!/usr/bin/env bash

while true; do
	nitrogen --set-zoom-fill --random /mnt/Hard320/reserve1/oboi/ --save 
	sleep 100
done
