# .config/sv/corectrl/run
#!/usr/bin/env bash

exec corectrl -t 540000
