# .config/sv/i3helper/run
#!/usr/bin/env bash

exec ~/.local/bin/I3Helper 13131
