# .config/sv/runsvdir-graphic/run
#!/bin/bash

svdir="$HOME/.local/share/service-graphic"

exec runsvdir "$svdir"
