# .config/sv/keyd-app/run
#!/usr/bin/env bash
exec keyd-application-mapper -v
